# logger
 #! /bin/sh
# /etc/init.d/logger

### BEGIN INIT INFO
# Provides:          logger
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Simple script to start a program at boot
# Description:       A simple script hutdown.
### END INIT INFO

# If you want a command to always run, put it here

# Carry out specific functions when asked to by the system
case "$1" in
  start)
    echo "Starting logger"
    # run application you want to start
    /home/pi/logger.sh
    ;;
  stop)
    echo "Stopping logger"
    # kill application you want to stop
    killall logger.sh
    killall gpxlogger
    ;;
  *)
    echo "Usage: /etc/init.d/logger {start|stop}"
    exit 1
    ;;
esac

exit 0 
